# private_dot_config/systemd/user/keyd-application-mapper.service
[Unit]
Description=Keyd Application Mapper
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/keyd-application-mapper -d
Restart=always

[Install]
WantedBy=default.target
